# patches/owner-console-v3.0.0/drivelab-owner-ops.service
[Unit]
Description=Execute allowlisted DriveLab Owner Console operations
After=local-fs.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/local/sbin/drivelab-owner-ops-worker
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
LockPersonality=true
UMask=0027
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
